# curiosum.service
[Unit]
Description=curiosum

[Service]
ExecStart=/home/eagle/curiosum/start.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target 
